# detections/local_rules.xml
This XML file does not appear to have any style information associated with it. The document tree is shown below.
<!-- 
  homelab-security custom rules
  ============================================================================
  ID range: 100100, 100200-100299 (per Wazuh custom-rule convention, > 100000)
  
  All <if_sid> / <if_group> values verified against the installed Wazuh 4.14.5
  ruleset at /var/ossec/ruleset/rules/. Stock-coverage analysis completed
  before authoring these rules — see ../docs/findings.md finding #10 for the
  redundancy mapping that resulted in 2 rules being removed from an earlier
  draft (Defender tampering → covered by stock 92008-92015; LSASS access →
  covered by stock 92900).

  MITRE ATT&CK tactics covered:
    TA0007 Discovery            — 100200
    TA0002 Execution            — 100210, 100250, 100290
    TA0005 Defense Evasion      — 100210, 100250, 100280
    TA0006 Credential Access    — 100240
    TA0003 Persistence          — 100260, 100270
    TA0011 Command & Control    — 100290
    TA0040 Impact               — 100280
    TA0008 Lateral Movement     — 100240 (adjacent)

  See ./README.md for per-rule writeups, validation status, and test methodology.
 -->
<group name="homelab,custom,">
<!--  =====================================================================
       100100 — UDR7 firewall block base rule (silent, level 0)
       Decoder: unifi-fw (custom, see decoders/unifi-fw.xml)
       Purpose: Fires on every block event from the UDR7 so 100200 can
                correlate on it via <if_matched_sid>. Does NOT generate
                alerts by itself (level 0).
       =====================================================================  -->
<rule id="100100" level="0">
<decoded_as>unifi-fw</decoded_as>
<field name="action">block</field>
<description>UDR7 firewall block (base rule, no alert).</description>
<group>unifi-fw,firewall,</group>
</rule>
<!--  =====================================================================
       100200 — Internal scan/probe (Discovery / T1046)
       Parent: 100100 (custom base rule, fires on every UDR7 firewall block)
       Logic: 20+ blocks from same source in 60s
       Note: Stock Wazuh does not ship a <different_dstport /> correlation
             modifier (only <different_url /> exists in the stock ruleset),
             so we use volume-only correlation from the same source. Less
             specific than per-port-diversity correlation but works with
             stock Wazuh's capabilities.
       =====================================================================  -->
<rule id="100200" level="10" frequency="20" timeframe="60">
<if_matched_sid>100100</if_matched_sid>
<same_source_ip/>
<description>Possible internal scan/probe: 20+ firewall blocks from $(srcip) in 60s.</description>
<mitre>
<id>T1046</id>
</mitre>
<group>recon,discovery,attack.discovery,</group>
</rule>
<!--  =====================================================================
       100210 — Encoded PowerShell command (Execution / Defense Evasion)
       Parent: sysmon_event1 (group — Sysmon EventID 1 process create)
       Verified: stock parent group is sysmon_event1 (NOT a specific rule ID)
       =====================================================================  -->
<rule id="100210" level="12">
<if_group>sysmon_event1</if_group>
<field name="win.eventdata.image" type="pcre2">(?i)\\powershell(_ise)?\.exe$</field>
<field name="win.eventdata.commandLine" type="pcre2">(?i)\s-e(c|nc|ncodedcommand)?\s+[A-Za-z0-9+/=]{40,}</field>
<description>PowerShell launched with -EncodedCommand (base64 payload).</description>
<mitre>
<id>T1059.001</id>
<id>T1027</id>
</mitre>
<group>powershell,obfuscation,attack.execution,attack.defense_evasion,</group>
</rule>
<!--  =====================================================================
       100240 — Internal SSH brute force (Credential Access / Lateral Movement)
       Parent: 5712 (stock sshd brute force detection)
       Adds: RFC1918 source filter to escalate severity for lateral attempts
       =====================================================================  -->
<rule id="100240" level="12">
<if_sid>5712</if_sid>
<srcip>192.168.0.0/16</srcip>
<description>SSH brute force from internal source $(srcip) — lateral movement candidate.</description>
<mitre>
<id>T1110.001</id>
<id>T1021.004</id>
</mitre>
<group>authentication_failures,attack.credential_access,attack.lateral_movement,</group>
</rule>
<!--  =====================================================================
       100250 — LOLBin abuse (Defense Evasion / Execution)
       Parent: sysmon_event1 group
       NOTE: certutil patterns removed — stock rules 92016-92018 already cover
             certutil masquerade and -decode. This rule covers mshta, regsvr32
             (Squiblydoo), rundll32 (javascript:), wmic xsl, bitsadmin — paths
             stock Wazuh does not detect.
       =====================================================================  -->
<rule id="100250" level="12">
<if_group>sysmon_event1</if_group>
<field name="win.eventdata.commandLine" type="pcre2">(?i)(mshta.*\b(http|javascript:)|regsvr32.*\/i:http|regsvr32.*scrobj\.dll|rundll32.*javascript:|wmic.*\/format:.*http|bitsadmin.*\/transfer)</field>
<description>LOLBin abuse pattern: $(win.eventdata.image)</description>
<mitre>
<id>T1218</id>
<id>T1105</id>
</mitre>
<group>lolbin,living_off_the_land,attack.defense_evasion,attack.execution,</group>
</rule>
<!--  =====================================================================
       100260 — Suspicious scheduled task (Persistence / T1053.005)
       Parent: sysmon_event1 group
       =====================================================================  -->
<rule id="100260" level="11">
<if_group>sysmon_event1</if_group>
<field name="win.eventdata.image" type="pcre2">(?i)\\schtasks\.exe$</field>
<field name="win.eventdata.commandLine" type="pcre2">(?i)\/create.*(\/sc\s+(onlogon|onstart)|\/ru\s+(system|highest)|\/tr.*(powershell|cmd|wscript|cscript|rundll32))</field>
<description>Suspicious scheduled task created.</description>
<mitre>
<id>T1053.005</id>
</mitre>
<group>persistence,scheduled_task,attack.persistence,</group>
</rule>
<!--  =====================================================================
       100270 — New Windows service (Persistence / T1543.003)
       Parent: 61138 (stock: "New Windows Service Created", Security 7045)
       
       Rationale: An earlier draft of this rule chained off sysmon_event1
       against sc.exe, which missed any service created via PowerShell
       New-Service or direct SCM API calls. Tiering on 61138 (which fires on
       Security channel EventID 7045) catches every service install path.
       
       Stock 61138 is level 5 (mostly informational because Windows installers
       create services routinely). This rule promotes to level 10 to surface
       for triage. Future tuning will add field filters on imagePath patterns
       (temp paths, script interpreters) once field-name verification is done
       against a live alert.
       =====================================================================  -->
<rule id="100270" level="10">
<if_sid>61138</if_sid>
<description>New Windows service installed (Security 7045) — triage required.</description>
<mitre>
<id>T1543.003</id>
</mitre>
<group>persistence,new_service,attack.persistence,</group>
</rule>
<!--  =====================================================================
       100280 — Critical service stop (Defense Evasion / Impact)
       Parent: 60002 (stock: System channel base rule)
       Match: EventID 7036 + message indicates "entered the stopped state"
              for a service on the security-relevant watchlist.
       
       Stock Wazuh ships rules for SOME service stops (e.g. Software Protection
       in 0585-win-application_rules.xml) but no generic SCM 7036 watchlist
       rule. This is a genuine gap.
       
       Watchlist intentionally tight: stopping any of these is rarely benign
       on a monitored endpoint.
       =====================================================================  -->
<rule id="100280" level="12">
<if_sid>60002</if_sid>
<field name="win.system.eventID">^7036$</field>
<field name="win.system.message" type="pcre2">(?i)(WinDefend|Sense|WdNisSvc|MpsSvc|EventLog|Sysmon64|WazuhSvc).{0,80}entered the stopped state</field>
<description>Security-relevant service stopped (EventID 7036).</description>
<mitre>
<id>T1562.001</id>
<id>T1489</id>
</mitre>
<group>defense_evasion,service_stop,attack.defense_evasion,attack.impact,</group>
</rule>
<!--  =====================================================================
       100290 — PowerShell download cradle (Ingress Tool Transfer / Execution)
       Parent: sysmon_event1 group
       =====================================================================  -->
<rule id="100290" level="13">
<if_group>sysmon_event1</if_group>
<field name="win.eventdata.image" type="pcre2">(?i)\\powershell(_ise)?\.exe$</field>
<field name="win.eventdata.commandLine" type="pcre2">(?i)(iex|invoke-expression).{0,80}(new-object.{0,40}(net\.webclient|system\.net\.webclient)|invoke-webrequest|iwr|wget|curl|downloadstring|downloaddata|downloadfile)</field>
<description>PowerShell download cradle.</description>
<mitre>
<id>T1105</id>
<id>T1059.001</id>
</mitre>
<group>ingress_tool_transfer,powershell,attack.command_and_control,attack.execution,</group>
</rule>
</group>
